# crontab
SHELL=/bin/bash
* * * * * root /usr/local/bin/php /your/path/to/oxwall/ow_cron/run.php
