# crontab
SHELL=/bin/bash
* * * * * root /usr/local/bin/php /var/www/html/ow_cron/run.php
